# build/tools/cloudcore.service
[Unit]
Description=cloudcore.service

[Service]
Type=simple
ExecStart=/etc/kubeedge/cloudcore

[Install]
WantedBy=multi-user.target
