# build/tools/edge.service
[Unit]
Description=edge.service

[Service]
Type=simple
ExecStart=/etc/kubeedge/edgecore

[Install]
WantedBy=multi-user.target
