# build/tools/edgecore.service
[Unit]
Description=edgecore.service

[Service]
Type=simple
ExecStart=/usr/local/bin/kubeedge/edgecore

[Install]
WantedBy=multi-user.target
